# host/linux-usb/rootfs-overlay/etc/systemd/system/getty@tty1.service.d/override.conf
[Unit]
# Disable tty1 getty to free the console for the auto-start service
ConditionPathExists=!/etc/systemd/system/likeos-autostart.service
